# src/Seagate/BA_NAS_4xx_GPL/opensource/rootfs/misc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
#
# The <id> field is used to specify the controlling TTY for the process.
#       This is different then sysvinit behavior!
#

# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
#null::sysinit:/bin/mount -o remount,rw /
null::sysinit:/bin/mount -n -t proc proc /proc
null::sysinit:/bin/mount -n -t sysfs none /sys
null::sysinit:/bin/mount -n -t usbfs /proc/bus/usb/ /proc/bus/usb/
null::sysinit:/bin/mount -n -t devpts /dev/pts /dev/pts
null::sysinit:/bin/mount -n -t securityfs securityfs /sys/kernel/security/
null::sysinit:/bin/mount -n -t tmpfs tmpfs /dev/shm -o size=512K
null::sysinit:/bin/mount --bind /dev/shm /etc/cups/certs
#null::sysinit:/bin/mount --bind /dev/shm /usr/var/locks
#null::sysinit:/bin/mount --bind /dev/shm /etc/cups
#null::sysinit:/bin/cp -a /etc/cups.default/* /etc/cups/.
null::sysinit:/sbin/swapon /dev/md1
null::sysinit:/bin/mount -o noatime /dev/md2 /atsahs
#null::sysinit:/bin/mount -a
null::sysinit:/bin/hostname -F /etc/hostname
#null::sysinit:/sbin/ifconfig lo 127.0.0.1 up
#null::sysinit:/sbin/route add -net 127.0.0.0 netmask 255.0.0.0 lo

# now run any rc scripts
# System startup
::sysinit:/etc/init.d/rcS
::sysinit:/sbin/insmod /lib/modules/2.6.12.6-arm1/usblp.ko
::respawn:/sbin/getty -L ttyS0 115200 vt100

# Set up a couple of getty's
#tty1::respawn:/sbin/getty 38400 tty1
#tty2::respawn:/sbin/getty 38400 tty2

# Put a getty on the serial port
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Logging junk
null::sysinit:/bin/touch /var/log/messages
#null::respawn:/sbin/syslogd -n -m 0
#null::respawn:/sbin/klogd -n
#tty3::respawn:/usr/bin/tail -f /var/log/messages

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
#null::shutdown:/usr/bin/killall klogd
#null::shutdown:/usr/bin/killall syslogd
#null::shutdown:/bin/umount -a -r
#null::shutdown:/sbin/swapoff -a

::sysinit:/etc/init.d/rcS

# Console

# Runtime
#::respawn:

# Go into runlevel 0 (Halt)
::shutdown:/etc/init.d/rcH

